# mtr-monitor/rootfs/etc/services.d/mtr-monitor/run
#!/command/execlineb -P
# Redirect stderr to stdout
2>&1
# Run the MTR monitoring script
exec /usr/local/bin/mtr-monitor.sh